# systemd/user/velyx-firstboot.service
[Unit]
Description=Velyx OS First Boot Entry
PartOf=velyx-session.target
After=velyx-permissions.service velyx-launcher.service velyx-session-manager.service graphical-session-pre.target dbus.service
Requires=velyx-permissions.service velyx-launcher.service velyx-session-manager.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=-%h/.config/velyx/velyx.env
Environment=VELYX_USER_ID=%u
Environment=HOME=%h
ExecStart=/usr/bin/velyx-entry --session
Restart=always
RestartSec=1
